# system/arch-notes/post-install
# swap file
su
cd /root
dd if=/dev/zero of=/swapfile bs=1M count=2048 status=progress
chmod 600 /swapfile
mkswap /swapfile

cp /etc/fstab /etc/fstab.bak
echo '/swapfile none swap sw 0 0' | tee -a /etc/fstab
swapon -a

# timedatectl list-timezones
timedatectl set-timezone Europe/Warsaw
systemctl enable systemd-timesyncd
hostnamectl set-hostname machine
nvim /etc/hosts
# add
# 127.0.0.1 localhost
# 127.0.1.1 machine

# install micro-code for CPU
pacman -S intel-ucode

# video driver, for Intel or AMD GPU is:
pacman -S mesa
# pacman -S nvidia
# pacman -S nvidia-lts
